# packaging/mwb.service
[Unit]
Description=Mouse Without Borders Linux client (bidirectional)
After=network-online.target graphical.target
Wants=network-online.target

[Service]
Type=simple
# DISPLAY and XAUTHORITY are auto-detected by the mwb binary at startup.
# Override here only if auto-detection fails for your setup.
# Environment=DISPLAY=:1
# Environment=XAUTHORITY=/run/user/1000/gdg/Xauthority
ExecStart=/usr/local/bin/mwb -bidi -edge left
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
